# edms/templates/mayan-uwsgi.conf.j2
[program:mayan-uwsgi]
command = {{ venv }}/bin/uwsgi --ini {{ venv }}/uwsgi.ini
user = root
autostart = true
autorestart = true
redirect_stderr = true
